# rootfs/etc/services.d/nginx/run
#!/command/execlineb -P
with-contenv
foreground { s6-echo "INFO: starting Nginx..." }
fdmove -c 2 1
nginx -g "daemon off;"