# rootfs/etc/services.d/ttyd/run
#!/command/execlineb -P
with-contenv
importas -S -D /home/jovyan SERVE_DIR
importas -S -D / NB_PREFIX
execline-cd ${SERVE_DIR}
foreground { s6-rmrf /var/run/shared/ttyd.sock }
backtick -n GIDLIST { id -G }
s6-envuidgid -B -i jovyan:shared
s6-applyuidgid -z -U
export HOME /home/jovyan
backtick -n CLEANED_NB_PREFIX { pipeline { echo ${NB_PREFIX} } sed -E "s|/+|/|g; s|/$||" }
importas -S -D "" CLEANED_NB_PREFIX
foreground { s6-echo "INFO: starting ttyd-server..." }
fdmove -c 2 1
ttyd --port 8886 --interface /var/run/shared/ttyd.sock --cwd ${SERVE_DIR} --url-arg --base-path ${CLEANED_NB_PREFIX}/ttyd /bin/bash