# rootfs/etc/services.d/notebook/run
#!/command/execlineb -P
with-contenv

# Import paths and base variables
importas -S -D /home/jovyan SERVE_DIRECTORY
importas -S -D /            NB_PREFIX

execline-cd ${SERVE_DIRECTORY}
foreground { s6-rmrf /var/run/shared/notebook.sock }

# Step 1: Capture PID 1 / Docker Host-provided supplementary groups (e.g., via --group-add)
backtick -n -E HOST_GIDS { id -G }

# Step 2: Capture the native groups assigned to user jovyan in the system database (/etc/group)
backtick -n -E DB_GIDS { id -G jovyan }

# Step 3: Dynamically look up the UID of jovyan and the GID of the requested PRIMARY_GID
backtick -n -E TARGET_UID { id -u jovyan }
backtick -n -E TARGET_GID { id -g jovyan }

# Step 4: Apply the identity configuration explicitly
# -u and -g set the execution context. -G applies the combined list of supplementary groups.
s6-applyuidgid -u ${TARGET_UID} -g ${TARGET_GID} -G "${HOST_GIDS} ${DB_GIDS}"

# Set user context environments
export HOME ${SERVE_DIRECTORY}

# Format and clean up the URL base path prefix
backtick -n -E -D "" CLEANED_NB_PREFIX { pipeline { echo ${NB_PREFIX} } sed -E "s|/+|/|g; s|/$||" }

foreground { s6-echo "INFO: starting notebook-server..." }
fdmove -c 2 1

# Execute Jupyter using uv
/opt/venv/bin/jupyter server --ServerApp.sock=/var/run/shared/notebook.sock --ServerApp.base_url="${CLEANED_NB_PREFIX}/notebook/"